# dnscrypt_dnsmasq/supervisord.conf
[supervisord]
nodaemon=true


[program:dnscrypt-proxy-serv]
command=/usr/sbin/dnscrypt-proxy -a 127.0.0.1:5300 -R opendns

[program:dnsmasq-serv]
command=/bin/bash -c "usr/sbin/dnsmasq -d"
